# tidbyt-web.service
[Unit]
Description=Tidbyt Matrix Display Service (Web Dashboard)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tidbyt
ExecStart=/usr/bin/python3 /home/pi/tidbyt/tidbyt_web.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
